# deploy/warsaw-beer-bot.service
[Unit]
Description=Warsaw Beer Crawler Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/warsaw-beer-bot
EnvironmentFile=/etc/warsaw-beer-bot/.env
ExecStart=/usr/bin/node /opt/warsaw-beer-bot/dist/index.js
Restart=on-failure
RestartSec=5
# Bot's signal handler stops crons + bot + db then exits, well under this
# budget. Cap below systemd's 90s default so a stuck process doesn't drag
# every redeploy.
TimeoutStopSec=20
StandardOutput=journal
StandardError=journal
User=warsaw-beer-bot
Group=warsaw-beer-bot

[Install]
WantedBy=multi-user.target
